# scripts/Linux/systemd/paper_playercount.service
[Unit]
Description=Run check_playercount script

[Service]
Type=oneshot
ExecStart=/usr/local/paper_server/check_playercount.sh